# ds4-xboxdrv.service
[Unit]
Description=DualShock 4/3/5 to Xbox 360 Controller Mapper
After=local-fs.target bluetooth.service
Requires=bluetooth.service

[Service]
Type=simple
User=root
Group=input
ExecStart=/usr/local/bin/ds4-xboxdrv.sh
Restart=on-failure
RestartSec=5s
TimeoutStopSec=15
KillMode=process
StartLimitIntervalSec=60
StartLimitBurst=5
StandardOutput=journal
StandardError=journal
ProtectSystem=strict
ProtectHome=true
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
