# pnp.service
[Unit]
Description=PlayStation to Xbox Controller Mapper
After=local-fs.target bluetooth.service
Requires=bluetooth.service

[Service]
Type=simple
ExecStart=/usr/bin/pnp-backend
Environment=WAYLAND_DISPLAY=wayland-0
Environment=XDG_RUNTIME_DIR=%t
Restart=on-failure
RestartSec=5s
TimeoutStopSec=15
KillMode=process
StandardOutput=journal
StandardError=journal

# Security Hardening
DynamicUser=yes
SupplementaryGroups=input
RuntimeDirectory=pnp
ConfigurationDirectory=pnp
ReadWritePaths=/dev/input /dev/uinput
ProtectSystem=strict
ProtectHome=true
NoNewPrivileges=true
RestrictRealtime=true
RestrictSUIDSGID=true
ProtectKernelModules=true
ProtectKernelTunables=true
ProtectControlGroups=true
MemoryDenyWriteExecute=true

[Install]
WantedBy=multi-user.target
